# deploy/ubuntu/systemd/cobien-update.service
[Unit]
Description=CoBien one-shot update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/cobien/cobien_FrontEnd/deploy/ubuntu
ExecStart=/bin/bash %h/cobien/cobien_FrontEnd/deploy/ubuntu/cobien-launcher.sh --mode update-once --non-interactive --yes --workspace %h/cobien --frontend-name cobien_FrontEnd --mqtt-name cobien_MQTT_Dictionnary --branch development_fix
